# scripts/systemd/utlab-pg-backup.timer
[Unit]
Description=Daily 04:30 KST UT Lab PostgreSQL backup

[Timer]
OnCalendar=*-*-* 04:30:00
Persistent=true

[Install]
WantedBy=timers.target
